# backup/dailygames-backup.timer
[Unit]
Description=Runs Daily Games Backup Service on a schedule

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
